# com.agenticedge.dailydigest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.agenticedge.dailydigest</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/pilli/agentic-digest/scripts/run.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/pilli/agentic-digest/logs/launchd_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/pilli/agentic-digest/logs/launchd_stderr.log</string>
	<key>KeepAlive</key>
	<false/>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
